# observability/fleet/systemd/fleet-dispatch-projection.timer
[Unit]
Description=Publish Edge dispatch-cycle projection every five minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=5m
Persistent=true

[Install]
WantedBy=timers.target
